# rabbitmqctl.ini
[program:rabbitmqctl]
command=/root/rabbitmqctl.sh
priority=3
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
