# trace.ini
[program:trace]
command=/usr/java/latest/bin/java -cp rabbittrace-0.1-jar-with-dependencies.jar edu.unc.tier.rabbittrace.Trace '#'
directory=/root
priority=2
